# 70-logitech-trackball.rules
ACTION=="add", ATTRS{name}=="Logitech USB Trackball", RUN+="/usr/local/bin/trackball
