# 70-trackball.rules
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c408", RUN+="/usr/local/bin/trackball.sh udev"